# Docker/cron/crontab
* * * * * echo "CRON IS RUNNING !!!" >> /partage/logs/activity.log
