# Docker/cron/supervisord.conf
[supervisord]
nodeamon=true

[program:cron]
command = cron -f
user = root
autostart = true